# ssh-scp/ansible-role-for-exporters/roles/exporter-linux/install-node-exporter/templates/node-exporter.service.j2
[Unit]
Description=Node Exporter
After=network.target

[Service]
Type=simple
User={{ node_exporter_user }}
Group={{ node_exporter_group }}
ExecStart={{ system_bin_dir }}/node_exporter
Restart=always

[Install]
WantedBy=multi-user.target
